# patch.1/etc/udev/rules.d/90-usb-device.rules
SUBSYSTEM!="usb", GOTO="usb_device_section_end"

RUN+="/usr/local/sbin/udev-usb.sh %E{ACTION} %N"

LABEL="usb_device_section_end"
