# patch.1/etc/xdg/autostart/caffeine.desktop
[Desktop Entry]
Name=Caffeine
Exec=caffeine.py
Type=Application
OnlyShowIn=MATE;
X-MATE-Autostart-Delay=10
