# patch.1/etc/xdg/autostart/mime-caja.desktop
[Desktop Entry]
Name=Mime Caja
Exec=xdg-mime default caja-folder-handler.desktop inode/directory
Type=Application
OnlyShowIn=MATE;
