# patch.1/etc/xdg/autostart/pa-vol-reset.desktop
[Desktop Entry]
Name=Reset volume
Exec=pa-vol.sh reset
Type=Application
OnlyShowIn=MATE;
X-MATE-Autostart-Delay=10
