# patch.1/etc/xdg/autostart/redshift.desktop
[Desktop Entry]
Name=Redshift
Exec=redshift-gtk
Type=Application
OnlyShowIn=MATE;
X-MATE-Autostart-Delay=10
