# patch.1/etc/profile.d/steam.sh
#!/bin/bash

export VERSION_ID=20200101000000
